# hesabixAPI/deployment/systemd/hesabix-rq-worker.service
[Unit]
Description=Hesabix RQ Worker (Background Jobs)
Documentation=https://hesabix.com/docs/rq-worker
After=network.target redis.service postgresql.service
Wants=redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/arc/hesabixAPI

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/var/www/arc/hesabixAPI"
Environment="PATH=/var/www/arc/hesabixAPI/.venv/bin:/usr/local/bin:/usr/bin:/bin"

# اجرای Worker
ExecStart=/var/www/arc/hesabixAPI/.venv/bin/python /var/www/arc/hesabixAPI/rq_worker.py

# Restart settings
Restart=always
RestartSec=10
StartLimitInterval=5min
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hesabix-rq-worker

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits - حذف محدودیت‌ها برای استفاده حداکثری از منابع
# MemoryMax: حذف شده - استفاده از RAM مورد نیاز
# CPUQuota: حذف شده - استفاده از CPU cores مورد نیاز
# Worker می‌تواند از منابع لازم برای پردازش سریع‌تر jobs استفاده کند

[Install]
WantedBy=multi-user.target

